# files/system/etc/profile.d/mizufetch.sh
alias neofetch='/usr/bin/mizufetch'
alias neowofetch='/usr/bin/mizufetch'
alias fastfetch='/usr/bin/mizufetch'
alias hyfetch='/usr/bin/mizufetch'
